# rootfs/etc/services.d/daemon/run
#!/usr/bin/with-contenv bash

s6-echo "Starting OpenHR20 daemon"
php-cgi /hr20/tools/daemon.php
