# rootfs/etc/services.d/daemon/run
#!/usr/bin/with-contenv bash

s6-echo "Starting OpenHR20 daemon"
/usr/bin/php /hr20/tools/daemon.php > /dev/null
